# ops/homeserver/addypin-backup.service
[Unit]
Description=addypin VPS → home server backup (pulselog --backup: DB + letsencrypt)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
# Run as the user whose ~/.ssh/addypin_vps key was authorized on the VPS.
# Override here or set via `systemctl edit` after install.
User=%i
EnvironmentFile=-/etc/default/addypin-backup
# pulselog --backup invokes addypin-pull.sh (the command source) to fetch the
# DB + letsencrypt into its staging dir, then tars/floors/publishes/rotates
# and emails on failure. A failed backup exits 1, surfaced by systemd.
# Adjust the binary path if your npm global prefix differs ($(command -v pulselog)).
ExecStart=/usr/local/bin/pulselog --backup --config /etc/addypin/pulselog.config.json
Nice=19
IOSchedulingClass=idle
