# ops/systemd/addypin-health.service
[Unit]
Description=addypin health check (pulselog: units, API, disk, cert, mail queue, DB)
After=network-online.target

[Service]
Type=oneshot
User=addypin
Group=addypin
# pulselog reads ops/pulselog/health.config.json, probes each check, stays
# silent on green, and emails the alert recipient via sendmail (Postfix →
# OpenDKIM) on any failure. The JSON config carries everything the bash
# script used to read from /etc/default/addypin-health, so no EnvironmentFile.
ExecStart=/usr/bin/node /opt/addypin/node_modules/pulselog/bin/pulselog.js --config /opt/addypin/ops/pulselog/health.config.json
Nice=10
