# ops/systemd/addypin-stats.service
[Unit]
Description=addypin daily stats snapshot (pin + customer counts to $DATA_DIR/stats.log)
After=network-online.target

[Service]
Type=oneshot
User=addypin
Group=addypin
WorkingDirectory=/opt/addypin
EnvironmentFile=/etc/addypin/env
# /bin/sh wrapper: systemd expands ${DATA_DIR} from EnvironmentFile, then sh
# does the >> redirect. RHEL 8 ships systemd 239 — StandardOutput=append:
# (240+) isn't available, so we go through sh. The script also reads
# DATA_DIR from its env to locate addypin.db, so the same env file pins
# both halves (db read path + log write path).
ExecStart=/bin/sh -c '/usr/bin/node --experimental-sqlite /opt/addypin/server/stats.js >> ${DATA_DIR}/stats.log'

# Match addypin.service hardening posture.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/addypin  # must match DATA_DIR in /etc/addypin/env
Nice=10
