# ops/systemd/addypin-stats-digest.timer
[Unit]
Description=Send the addypin weekly stats digest every Sunday morning

[Timer]
OnCalendar=Sun *-*-* 09:00:00
AccuracySec=15min
Persistent=true

[Install]
WantedBy=timers.target
